# youtube-transcriber.service
[Unit]
Description=YouTube Transcriber Service
After=network.target

[Service]
Type=simple
User=ytservice
Group=ytservice
WorkingDirectory=/opt/youtube-transcriber
Environment=PYTHONPATH=/opt/youtube-transcriber
ExecStart=/usr/local/bin/poetry run python server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target 